# audio-state-monitor/rootfs/etc/services.d/monitor/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Audio State Monitor stopped"
